# testjob.conf
description "Test job as template for XMage server"
author "David R. Smith"

start on runlevel [2345]

exec /usr/games/mage-server/startServer.sh
exec cpulimit -e java -l 50 -z
